# etc/systemd/minecraft-backup.timer
[Unit]
Description=run minecraft backup daily

[Timer]
OnCalendar=daily
Persistent=false

[Install]
WantedBy=timers.target